# templates/marathon.service.erb
[Unit]
Description=<%= @service_name %>
After=network.target


[Service]
<% @options.each do | option, value| -%>
<%= "Environment=\"MARATHON_#{option}=#{value}\"" %>
<% end -%>
User=root
Group=root
UMask=0027
SyslogIdentifier=<%= @service_name %>
WorkingDirectory=<%= @install_dir %>
ExecStart=<%= @install_dir %>/bin/start
LimitNOFILE=16384
CPUAccounting=true
MemoryAccounting=true
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target